# cpu-pinning.service
[Unit]
Description=Qubes CPU pinning
After=qubesd.service

[Service]
ExecStart=/usr/local/bin/cpu-pinning.py

[Install]
WantedBy=multi-user.target
